# redis.conf
description "redis server"

start on runlevel [23]
stop on shutdown
expect daemon
respawn

exec sudo -u redis redis-server /etc/redis/redis.conf
